# init.d/updateName
#!/bin/sh /etc/rc.common

START=10

BIN=/usr/bin/updateName

boot(){
  $BIN
}